# esp32-c3/mediaplayer-mqtt-translator.service
[Unit]
Description=Oled mediaplayer mqtt translator
After=network.target

[Service]
Type=simple
User=dietpi
Group=dietpi
WorkingDirectory=/home/dietpi/oled-mediaplayer
EnvironmentFile=/home/dietpi/oled-mediaplayer/.env
ExecStart=/usr/bin/python3 /home/dietpi/oled-mediaplayer/esp32-c3/mqtt_cover_translator.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
